# itservicedesk.service
[Unit]
Description=Calpion IT Service Desk
After=network.target postgresql.service

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/itservicedesk
ExecStart=/usr/bin/node server-production.cjs
Environment=NODE_ENV=production
Environment=PORT=3000
Restart=always
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=itservicedesk

[Install]
WantedBy=multi-user.target
