# servicedesk.service
[Unit]
Description=IT Service Desk Application
Documentation=https://github.com/your-repo/servicedesk
After=network.target postgresql.service

[Service]
Type=simple
User=servicedesk
Group=servicedesk
WorkingDirectory=/opt/servicedesk
Environment=NODE_ENV=production
Environment=PORT=5000
EnvironmentFile=/opt/servicedesk/.env
ExecStart=/usr/bin/npm run start
ExecReload=/bin/kill -HUP $MAINPID
KillMode=mixed
KillSignal=SIGINT
TimeoutStopSec=5
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=servicedesk

# Security settings
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ReadWritePaths=/opt/servicedesk/logs
ProtectHome=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes

[Install]
WantedBy=multi-user.target